# project/scripts/package/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/gpio-mqtt

update-rc.d gpio-mqtt defaults

#service gpio-mqtt start